# Hexasense_API/hexaboot.service
[Unit]
After=network.target

[Service]
ExecStart=/home/team-member/Codebase/Experiments/fastapi/hexaboot.sh

[Install]
WantedBy=default.target
