# Hexasense_API/hexaboot.service
[Unit]
After=network.target

[Service]
ExecStart=/etc/systemd/system/hexaboot.sh

[Install]
WantedBy=default.target
